# deploy/systemd/structai-worker.service
[Unit]
Description=StructAI arq import worker
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target
# The orchestrator and API share state via Postgres + Redis. Worker can
# come up before or after the API; if started together they coordinate
# via the queue.
PartOf=structai.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/structai-v2/backend
EnvironmentFile=/etc/structai.env
# Invoke the venv's arq directly. See structai-api.service for why we
# avoid `uv run` under ProtectHome=read-only.
ExecStart=/home/ubuntu/structai-v2/backend/.venv/bin/arq structai.worker.main.WorkerSettings
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/structai-v2 /var/lib/structai

[Install]
WantedBy=multi-user.target
