# deploy/systemd/structai-api.service
[Unit]
Description=StructAI FastAPI backend
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/structai-v2/backend
EnvironmentFile=/etc/structai.env
# Use uv to invoke uvicorn inside the project virtualenv. PATH ensures uv
# (installed via the user-mode astral installer) is reachable from systemd's
# minimal environment.
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ubuntu/.local/bin/uv run uvicorn structai.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/structai-v2 /var/lib/structai

[Install]
WantedBy=multi-user.target
