# deploy/systemd/structai-api.service
[Unit]
Description=StructAI FastAPI backend
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/structai-v2/backend
EnvironmentFile=/etc/structai.env
# Invoke the venv's uvicorn directly. `uv run` is avoided here because it
# tries to write to ~/.cache/uv at startup and ProtectHome=read-only
# blocks that. deploy.sh builds the venv with `uv sync` ahead of time.
ExecStart=/home/ubuntu/structai-v2/backend/.venv/bin/uvicorn structai.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/structai-v2 /var/lib/structai

[Install]
WantedBy=multi-user.target
